# docker/novnc/conf.d/xterm.conf
[program:xterm]
command=xterm
autorestart=false
